# systemd/rvc2hass.service
[Unit]
Description=RV-C CAN Bus to Home Assistant via MQTT
After=network.target can0-up.service
Wants=network.target
Requires=can0-up.service

[Service]
Type=simple
User=pi
ExecStart=/opt/rvc2hass/.venv/bin/python -m rvc2hass --profile /opt/rvc2hass/profiles/thor_hurricane_35m.yaml
WorkingDirectory=/opt/rvc2hass
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
